# packages/autostop-ec2/scripts/autostop.service
[Unit]
Description=AutoStop

[Service]
Type=simple
ExecStart=/usr/local/bin/autostop.sh

[Install]
WantedBy=multi-user.target
